# zapret2/service/zapret2.service
[Unit]
Description=Zapret2 v7.0 DPI Bypass Service
After=network-online.target
Wants=network-online.target

# 防止重复启动
StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/root/catmi/Zapret2/bin/zapret2d
ExecStop=/bin/kill -TERM $MAINPID

# 崩溃自动重启
Restart=always
RestartSec=2

# Watchdog（zapret2d 会自动发送心跳）
WatchdogSec=10
NotifyAccess=all

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# ================================
# 安全沙箱（生产级）
# ================================

# 文件系统保护
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# 禁止访问 /dev/shm、/dev/mqueue
PrivateIPC=yes

# 禁止提权
NoNewPrivileges=yes

# 限制网络族（只允许 IPv4/IPv6/NFQUEUE）
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK

# 限制系统调用（最小化攻击面）
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @raw-io @reboot

# 限制能力（只保留 CAP_NET_ADMIN）
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN

# 只允许写入日志目录
ReadWritePaths=/var/log/zapret2

# 禁止访问用户目录
InaccessiblePaths=/home
InaccessiblePaths=/root
InaccessiblePaths=/srv

# 禁止访问设备
DevicePolicy=closed

# ================================
# 环境变量
# ================================
Environment=ZAPRET2_DIR=/root/catmi/Zapret2

[Install]
WantedBy=multi-user.target
