# zapret2/service/zapret2.service
[Unit]
Description=Zapret2 v7.0 DPI Bypass Service
After=network-online.target
Wants=network-online.target

StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/root/catmi/Zapret2/bin/zapret2d
WorkingDirectory=/root/catmi/Zapret2
ExecStop=/bin/kill -TERM $MAINPID

Restart=always
RestartSec=2

WatchdogSec=10
NotifyAccess=all

LimitNOFILE=65535
LimitNPROC=4096

# ================================
# 安全沙箱（优化后）
# ================================

# 允许 zapret2d 访问自身目录（否则 226）
ProtectSystem=full
ReadWritePaths=/root/catmi/Zapret2
ReadWritePaths=/var/log/zapret2

# 允许访问 root 目录（否则无法读取配置）
ProtectHome=read-only
InaccessiblePaths=/home
InaccessiblePaths=/srv

# 允许访问 /tmp（nfqws2 需要）
PrivateTmp=no

# 允许访问 tun、netlink、proc
PrivateDevices=no
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# IPC 不能隔离，否则 NFQUEUE 失败
PrivateIPC=no

# 允许提权（否则 CAP_NET_ADMIN 无效）
NoNewPrivileges=no

# 网络族限制（保持安全）
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK

# 系统调用过滤（保留 NFQUEUE 所需）
SystemCallFilter=@system-service
SystemCallFilter=~@mount @reboot

# 能力（NFQUEUE + iptables 必须）
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

DevicePolicy=closed

Environment=ZAPRET2_DIR=/root/catmi/Zapret2

[Install]
WantedBy=multi-user.target
