# zapret2/service/zapret2.service
[Unit]
Description=Zapret2 v7.0 DPI Bypass Service
After=network-online.target
Wants=network-online.target

StartLimitIntervalSec=30
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/root/catmi/Zapret2/bin/zapret2d
WorkingDirectory=/root/catmi/Zapret2
ExecStop=/bin/kill -TERM $MAINPID

Restart=always
RestartSec=2

# ❌ zapret2d 不支持 systemd-notify，必须关闭
WatchdogSec=0
NotifyAccess=none

LimitNOFILE=65535
LimitNPROC=4096

# ================================
# 安全沙箱（兼容 zapret2d）
# ================================
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=no
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

PrivateIPC=no
NoNewPrivileges=no

RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK

SystemCallFilter=@system-service
SystemCallFilter=~@mount @reboot

CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

ReadWritePaths=/root/catmi/Zapret2
ReadWritePaths=/var/log/zapret2

InaccessiblePaths=/home
InaccessiblePaths=/srv

DevicePolicy=closed

Environment=ZAPRET2_DIR=/root/catmi/Zapret2

[Install]
WantedBy=multi-user.target
